# deploy/macos/launchd/com.voxbot.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.voxbot.bot</string>
  <key>UserName</key><string>voxbot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/voxbot/apps/voxbot/run-bot.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>/Users/voxbot/logs/bot.out.log</string>
  <key>StandardErrorPath</key><string>/Users/voxbot/logs/bot.err.log</string>
</dict>
</plist>
